# lms-wanguard-update-db.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root

*/30 * * * * root php /var/www/html/lms/plugins/LMSWanguardPlugin/bin/lms-wanguard-update-db.php
